# cmd/nitrod/nitrod.service
[Unit]
Description=nitrod-api
Documentation=https://getnitro.sh
Wants=network-pre.target
After=network-pre.target
StartLimitIntervalSec=0
StartLimitBurst=0
MemoryMax=256M

[Service]
ExecStart=/usr/sbin/nitrod
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
