# .docker/supervisor.conf
[supervisord]
nodaemon=true

[program:caddy]
command=/usr/bin/caddy run --resume

[program:nitrod]
command=/usr/bin/nitrod